# example-deactivate-check.service
[Unit]
Description=Check if unit example.service should be stopped.

[Service]
Type=oneshot
ExecStart=/absolute/path/to/example-check-stop.sh
